# stage5/00-rasppleii-tweaks/files/rasppleii_boot.service
[Unit]
Description=Misc boot actions for Raspple II

[Service]
Type=oneshot
ExecStart=/etc/rc.rasppleii_boot

[Install]
WantedBy=multi-user.target
